# debian/postinst
#!/bin/bash
systemctl enable --now qomui.service