# debian/prerm
#!/bin/bash
systemctl disable --now qomui.service